# weston.service
[Unit]
Description=Weston Compositor

[Service]
ExecStart=/usr/bin/weston

[Install]
WantedBy=graphical-session.target
